# files/redhat/kong.service
[Unit]
Description=Kong

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/kong start
ExecStop=/usr/local/bin/kong stop

[Install]
WantedBy=multi-user.target
